# deploy/systemd/second-lane.service
[Unit]
Description=Second Lane
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/second-lane
EnvironmentFile=/opt/second-lane/.env
ExecStart=/opt/second-lane/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8787
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
